# deploy/systemd/whatsapp-chat-system.service
[Unit]
Description=WhatsApp Chat System standalone API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/whatsapp-chat-system
EnvironmentFile=-/etc/whatsapp-chat-system/api.env
# Required from api.env: DATABASE_URL and WHATSAPP_BRIDGE_INTERNAL_TOKEN.
Environment=CHAT_SYSTEM_RUNTIME_DIR=/var/lib/whatsapp-chat-system/api
Environment=CHAT_SYSTEM_WEB_DIST=/opt/whatsapp-chat-system/web/dist
ExecStart=/opt/whatsapp-chat-system/.venv/bin/python -m whatsapp_chat_system.cli serve --host 127.0.0.1 --port 8792 --web-dist /opt/whatsapp-chat-system/web/dist
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
